# deploy/clibot.service
# clibot systemd user service file
# Install to: ~/.config/systemd/user/clibot.service
# Then run: systemctl --user daemon-reload && systemctl --user enable clibot && systemctl --user start clibot
#
# === CUSTOMIZABLE PATHS ===
# You can change these paths according to your setup:
#   - ExecStart:                Binary path and config file location
#   - WorkingDirectory:         Working directory (default: user's home directory)
#   - StandardOutput/Error:     Logging configuration (using journal by default)
#
# Binary location:
#   - After `go install`:        ~/go/bin/clibot (if GOPATH=~/go)
#   - Or if in PATH:             clibot
#
# Config location (MUST specify --config):
#   - User-specific:  ~/.config/clibot/config.yaml (recommended)
#   - Project dir:     ./config.yaml
#   - Custom path:     /any/path/config.yaml

[Unit]
Description=Clibot - IM to AI CLI bridge
Documentation=https://github.com/keepmind9/clibot
After=network.target
Wants=network-online.target

[Service]
# Working directory (YOU CAN CUSTOMIZE)
# Comment out to use user's home directory, or specify a custom path
# WorkingDirectory=~/projects

# ExecStart - Binary path and config file (YOU CAN CUSTOMIZE BOTH)
# Default: Uses ~/go/bin/clibot (installed via 'go install')
# If ~/go/bin is in PATH, you can use 'clibot' instead
ExecStart=~/go/bin/clibot serve --config ~/.config/clibot/config.yaml

# Auto-restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clibot

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
