# deploy/clibot.conf
# clibot supervisor configuration
# Install to: /etc/supervisor/conf.d/clibot.conf
# Then run: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start clibot
#
# NOTE: supervisor is a system-level service and requires root to install/configure
# For user-level service management without root, consider using systemd user services:
#   cp deploy/clibot.service ~/.config/systemd/user/
#   systemctl --user enable clibot && systemctl --user start clibot
#
# === CUSTOMIZABLE PATHS ===
# You can change these paths according to your setup:
#   - command:     Binary path and config file location
#   - user:        User to run clibot (default: current user)
#   - stdout_log:  Log file location
#
# Common config locations:
#   - User-specific:  ~/.config/clibot/config.yaml (default, recommended)
#   - System-wide:     /etc/clibot/config.yaml (requires root)
#   - Project dir:     ./config.yaml
#   - Custom path:     /any/path/config.yaml
#
# Example: Use custom config location
#   command=/usr/local/bin/clibot serve --config /custom/path/config.yaml

[program:clibot]
; Command to run clibot - Binary path and config file (YOU CAN CUSTOMIZE BOTH)
; Default: Uses ~/go/bin/clibot with ~/.config/clibot/config.yaml
command=/home/yourusername/go/bin/clibot serve --config /home/yourusername/.config/clibot/config.yaml

; User to run clibot (YOU CAN CUSTOMIZE)
; Comment out to run as current user (recommended), or specify a dedicated user
; user=clibot

; Start automatically when supervisor starts
autostart=true

; Auto-restart if clibot crashes
autorestart=true

; Number of retries on startup
startretries=3

; Wait 10 seconds before restarting
startsecs=10

; Number of worker processes (1 for single instance)
numprocs=1

; Process name format
process_name=%(program_name)s

; Redirect stderr to stdout
redirect_stderr=true

; Log files (YOU CAN CUSTOMIZE)
; Default: Uses user-specific log location
stdout_logfile=~/clibot.log
; Or use system-wide log location:
; stdout_logfile=/var/log/clibot/clibot.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

; Environment variables (optional)
; environment=PATH="/usr/local/bin:/usr/bin:/bin"

; Priority (higher = starts first)
priority=100

; Stop signal
stopsignal=TERM

; Wait time for graceful shutdown
stopwaitsecs=30

; Run as daemon
daemon=false

; Kill as group
killasgroup=true
